# cron/crontab
@daily bash /usr/src/cron/backupTaigaDB
